# deploy/ride-optimizer.service
[Unit]
Description=Ride Optimizer — start container on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/admin/ride-optimizer
ExecStartPre=+/bin/chown -R admin:admin /home/admin/ride-optimizer/config
ExecStart=/home/admin/.local/bin/podman-compose up -d
ExecStop=/home/admin/.local/bin/podman-compose down
User=admin
Environment=PATH=/home/admin/.local/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ride-optimizer

[Install]
WantedBy=multi-user.target
